# etc/profile.d/go.sh
# Begin /etc/profile.d/go.sh

export PATH=/usr/local/go/bin:$PATH

# End /etc/profile.d/go.sh

