# etc/profile.d/kotlin.sh
# Begin /etc/profile.d/kotlin.sh

# Set JAVA_HOME directory
JAVA_HOME=

# Adjust PATH
export PATH=$PATH:/opt/kotlinc/bin:/opt/gradle-8.12/bin:/opt/android-sdk/cmdline-tools/latest/bin:/opt/android-sdk/emulator:/opt/android-sdk/platform-tools/
export ANDROID_HOME=/opt/android-sdk

# End /etc/profile.d/kotlin.sh
