# etc/profile.d/mapkg.sh
# Begin /etc/profile.d/mapkg.sh

export PATH=/opt/mapgk:$PATH
export PKG_CONFIG_PATH=/opt/mapkg/lib:$PKG_CONFIG_PATH
export CPLUS_INCLUDE_DIR=/opt/mapkg/include:$CPLUS_INCLUDE_DIR

# End /etc/profile.d/mapkg.sh
