# etc/profile.d/mapkg.sh
# Begin /etc/profile.d/mapkg.sh
export PATH=/opt/mapgk:/opt/mapkg/bin:$PATH
export PKG_CONFIG_PATH=/opt/mapgk/lib:$PKG_CONFIG_PATH
export GCC_INCLUDE_DIR=/opt/mapkg/include:$GCC_INCLUDE_DIR
# End /etc/profile.d/mapkg.sh
