# etc/profile.d/rustc.sh
# Begin /etc/profile.d/rustc.sh

export PATH=/opt/rustc-1.80.1/bin:$PATH

# End /etc/profile.d/rustc.sh

