# etc/profile.d/zig.sh
# Begin /etc/profile.d/zig.sh

export PATH=$HOME/opt/zig:$PATH

# End /etc/profile.d/zig.sh
